# .config/systemd/user/cli-update.service
[Unit]
Description=Update cli toolbox

[Service]
ExecStart=distrobox upgrade default-cli
