# .config/systemd/user/cli-update.service
j[Unit]
Description=Update cli toolbox

[Service]
ExecStart=/var/home/halli/.local/bin/update-distrobox.sh
