# .config/systemd/user/rustup.service
[Unit]
Description=Update rust

[Service]
ExecStart=/var/home/halli/.local/bin/update-rust.sh
