# .config/systemd/user/rustup.service
[Unit]
Description=Update rust

[Service]
ExecStart=rustup update
